# yay-clear/.config/systemd/user/yay-clear.service
[Unit]
Description=Clear yay cache

[Service]
Type=oneshot
ExecStart=/usr/bin/yay -Sc --aur --noconfirm
