# battery/.config/systemd/user/battery_notify.timer
[Unit]
Description=Check battery every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
