# battery/.config/systemd/user/battery_notify.service
[Unit]
Description=Battery level notification

[Service]
ExecStart=%h/.local/bin/battery_notify.sh
